# data/gamemoded.service.in
[Unit]
Description=gamemoded

[Service]
Type=dbus
BusName=com.feralinteractive.GameMode
ExecStart=@BINDIR@/gamemoded -l

[Install]
WantedBy=default.target
